# templates/default/systemd/sockd.service.erb
# Dynamically generated by Chef - local modifications will be overwritten

[Unit]
Description=Dante SOCKS v4 and v5 compatible proxy server and client
After=network.target remote-fs.target nss-lookup.target

[Service]
Type=forking
ExecStartPre=/bin/touch <%= node[:dante][:daemon][:pid_file] %>
ExecStartPre=/bin/chown <%= node[:dante][:auth][:users][:privileged] %> <%= node[:dante][:daemon][:pid_file] %>
ExecStart=<%= node[:dante][:daemon][:binary] %> -f <%= node[:dante][:configuration_file] %> -p <%= node[:dante][:daemon][:pid_file] %> -D

[Install]
WantedBy=multi-user.target
